# gpioviz.service
[Unit]
Description=GPIO Visualizer Web Application
After=network.target

[Service]
Type=simple
User=peter
WorkingDirectory=/home/peter/gpioviz
ExecStart=/home/peter/gpioviz/venv/bin/python /home/peter/gpioviz/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
